# jobs/build/ami-updater/roles/registry_credentials_sync/files/ops-registry-secret-sync.service
[Unit]
Before=atomic-openshift-node.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ops-registry-secret-sync.sh
SyslogIdentifier=ops-registry-secret-sync
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
